# meta-opendsp/recipes-core/opendspd/files/first-boot-setup.service
[Unit]
Description=First Boot Setup Script
DefaultDependencies=no
After=local-fs.target systemd-remount-fs.service var.mount
Before=sysinit.target shutdown.target
ConditionPathExists=!/var/lib/first-boot-setup/completed

[Service]
Type=oneshot
ExecStart=/usr/sbin/first-boot-setup.sh
RemainAfterExit=no
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=multi-user.target
