# PKGBUILDs/opendspd/systemd/opendsp.service
[Unit]
Description=OpenDSP Remote Control Service
After=remote-fs.service
#ConditionCapability=CAP_SYS_NICE

[Service]
Type=simple
User=opendsp
Group=opendsp
WorkingDirectory=/home/opendsp/
ExecStart=/usr/bin/opendspd
Environment="DISPLAY=:0"
Environment="XAUTHORITY=/tmp/.Xauthority"
StandardOutput=syslog
StandardError=syslog
#CPUSchedulingPolicy=rr
#CPUSchedulingPriority=50
LimitRTPRIO=infinity
LimitMEMLOCK=infinity
LimitRTTIME=infinity
TasksMax=infinity

[Install]
WantedBy=multi-user.target
