# PKGBUILDs/opendspd/systemd/startx.service
[Unit]
After=systemd-user-session.service

[Service]
ExecStart=/bin/su opendsp -l -c /usr/bin/startx -- VT08

[Install]
WantedBy=multi-user.target
